# stow/macos/Library/LaunchAgents/com.sogamoso.workhours.caffeinate-run.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sogamoso.workhours.caffeinate-run</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/caffeinate</string>
    <string>-imu</string>
    <string>-t</string>
    <string>39600</string>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/tmp/com.sogamoso.workhours.caffeinate-run.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.sogamoso.workhours.caffeinate-run.err</string>
</dict>
</plist>
