# stow/macos/Library/LaunchAgents/com.sogamoso.workhours.caffeinate-watch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sogamoso.workhours.caffeinate-watch</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>exec "$HOME/.config/workhours/workhours-caffeinate.sh"</string>
  </array>

  <key>StartInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/tmp/com.sogamoso.workhours.caffeinate-watch.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/com.sogamoso.workhours.caffeinate-watch.err</string>
</dict>
</plist>
